# backuper.cron
0 0 * */1 * /usr/local/emhttp/plugins/backuper/bin/startSchedule.php
